# guest/desktop/rootfs-overlay/etc/systemd/system/dory-first-boot.service
[Unit]
Description=Prepare the Dory desktop machine
After=local-fs.target systemd-remount-fs.service
Before=NetworkManager.service lightdm.service ssh.service dory-boot.service
ConditionPathExists=!/var/lib/dory/first-boot-complete

[Service]
Type=oneshot
ExecStart=/usr/lib/dory/first-boot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
